# pdf_converter_v2/scripts/mineru-api.service
# MinerU file_parse API - 最简 systemd 服务（NVIDIA GPU）
# 安装：sudo cp mineru-api.service /etc/systemd/system/
#       sudo systemctl daemon-reload && sudo systemctl enable --now mineru-api

[Unit]
Description=MinerU file_parse API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/workspeac/Clerk2.5
Environment=MINERU_MODEL_SOURCE=modelscope
# 不加载 vLLM 插件，避免 PaddleX GenAI 插件在纯 NVIDIA 环境下报错
Environment=VLLM_PLUGINS=
# 跳过模型源连通性检查，加快启动（可选）
Environment=PADDLE_PDX_DISABLE_MODEL_SOURCE_CHECK=True
ExecStart=/root/workspeac/Clerk2.5/.venv_paddleocr/bin/mineru-api --host 0.0.0.0 --port 5282
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
